# deploy/ferretdb/systemd/librechat-ferretdb-ship-backups.timer
[Unit]
Description=Run LibreChat backup shipping after daily backup

[Timer]
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=20m
Persistent=true

[Install]
WantedBy=timers.target
